# noagenda.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

MAILTO=donaghhatton@gmail.com
59 19,21,23 * * 4,7 cd /home/ubuntu/venv/noagenda-db/noagenda-db; ./fetch_updates.sh
