# deploy/nadb/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = nginx

[program:gunicorn]
command = gunicorn --bind localhost:8001 nadb.wsgi:application
